# isanteplus-monitor.service
[Unit]
Description=iSantePlus Server Monitor
After=network.target tomcat7.service mysql.service

[Service]
Type=simple
ExecStart=/usr/local/bin/isanteplus-monitor.sh
Restart=on-failure
RestartSec=10
NoNewPrivileges=yes
ProtectHome=yes
# Use 'full' instead of 'strict' — the monitor needs to call
# systemctl restart tomcat7, which requires write access to systemd paths.
ProtectSystem=full
ReadWritePaths=/var/log/isanteplus-monitor /tmp

[Install]
WantedBy=multi-user.target
